# deploy/launchd/com.pulse.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  pulse DB バックアップ(§8: 夜間 pg_dump を Mac へ)。

  radio バッチ(04:30)の後、Mac がまだ起きている 05:15 に Pi から
  pg_dump を pull し、episodes/ もミラーする。Mac が閉じていた日は
  スキップされるだけ(縮退許容。月次の実在確認は定常運用タスク)。
  CHANGEME を自分のユーザー名に置換してから使う(deploy/mac.md 10章)。
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pulse.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/CHANGEME/pulse/bin/backup-pulse-db.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/CHANGEME/pulse/logs/backup.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/CHANGEME/pulse/logs/backup.err.log</string>
</dict>
</plist>
